# lib/init/kiln-controller.service
[Unit]
Description=kiln-controller

[Service]
# __KILN_DIR__ and __UV__ are substituted with real paths by start-on-boot.
WorkingDirectory=__KILN_DIR__
ExecStart=__UV__ run kiln-controller.py

[Install]
WantedBy=multi-user.target
